# build/aero-teraranger.service
[Unit]
Description=Serial teraranger for intel aero

[Service]
Type=simple
ExecStart=/usr/bin/intel_tera_mavlink
Restart=on-failure

[Install]
WantedBy=multi-user.target
